# bingwallpaper/bingwallpaper.timer
[Unit]
Description=Run bing wallapper hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
