# bingwallpaper/bingwallpaper.timer
[Unit]
Description=Run bing wallapper daily

[Timer]
OnCalendar=daily
Persistent=true #Realtime timer

[Install]
WantedBy=timers.target
